# agent_team/launchd/com.marketdaily.agentteam.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marketdaily.agentteam.digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/delvin/Downloads/Delvin agent/agent_team/run_digest.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>3</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/delvin/Downloads/Delvin agent/agent_team/logs/launchd-digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/delvin/Downloads/Delvin agent/agent_team/logs/launchd-digest.err</string>
</dict>
</plist>
